# scripts/systemd/keipes-ai-mcp.service
[Unit]
Description=Keipes AI MCP Server
After=network.target

[Service]
Type=simple
User=ec2-user
Group=ec2-user
WorkingDirectory=/home/ec2-user
ExecStart=/usr/local/bin/keipes-ai-mcp
Restart=always
RestartSec=5
Environment=RUST_LOG=info
Environment=PORT=8080

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target